# deploy/systemd/signal-flow-pipeline.timer
[Unit]
Description=Run Signal-Flow pipeline every 4 hours

[Timer]
OnCalendar=*-*-* 00/4:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
